# build/sysroot/usr/lib/systemd/system/casaos-app-management.service
[Unit]
After=casaos-message-bus.service
After=docker.service
Description=CasaOS App Management Service

[Service]
ExecStartPre=/usr/bin/casaos-app-management -v
ExecStart=/usr/bin/casaos-app-management -c /etc/casaos/app-management.conf
PIDFile=/var/run/casaos/app-management.pid
Restart=always
Type=notify

[Install]
WantedBy=multi-user.target
